# eos-relocate-codecs.service.in
# Makes sure that codecs are present under /var/lib/codecs after boot.

[Unit]
Description=Helper script to deal with extra codecs
ConditionPathExists=/usr/share/EndlessOS/codecs

# Don't run 2 ldconfigs concurrently
After=ldconfig.service

[Service]
Type=oneshot
ExecStart=@bindir@/eos-relocate-codecs
User=root

[Install]
WantedBy=multi-user.target
