# eos-gstreamer-codecs-update.service.in
[Unit]
Description=Update or install codecs via GStreamer plugins

[Service]
Type=simple
ExecStart=@pkglibexecdir@/eos-gstreamer-codecs-update
