# eos-relocate-gst-plugins.service.in
# Makes sure that gstreamer plugins are present
# under /var/lib/codecs/gstreamer-1.0 after boot.

[Unit]
Description=Helper script to deal with extra codecs
ConditionPathExists=/usr/share/EndlessOS/gstreamer-1.0

[Service]
Type=oneshot
ExecStart=@bindir@/eos-relocate-gst-plugins
User=root

[Install]
WantedBy=multi-user.target
