# data/eos-codecs-activate.service.in
# Makes sure that codecs are present under /var/lib/codecs after boot.

[Unit]
Description=Helper script to deal with extra codecs
ConditionPathExists=/usr/share/eos-codecs-manager

# Needs to be in sysinit.target to run before systemd-update-done.service.
DefaultDependencies=no
Conflicts=shutdown.target
After=local-fs.target
Before=sysinit.target shutdown.target systemd-update-done.service

# Don't run 2 ldconfigs concurrently
After=ldconfig.service

[Service]
Type=oneshot
ExecStart=@pkglibexecdir@/eos-codecs-activate
User=root

[Install]
WantedBy=sysinit.target
