# data/eos-chromium-codecs-update.service.in
[Unit]
Description=Update or install Chromium codecs

[Service]
Type=simple
ExecStart=@pkglibexecdir@/eos-chromium-codecs-update
